# root/etc/s6-overlay/s6-rc.d/tcd/run
#!/command/with-contenv sh

# start daemon
exec /usr/local/bin/tcd ${URFD_CONFIG_DIR}/tcd.ini
