# root/etc/s6-overlay/s6-rc.d/urfd/run
#!/command/with-contenv sh

INIFILE="/config/urfd.ini"

# start daemon
exec /usr/local/bin/urfd ${INIFILE}
