# root/etc/s6-overlay/s6-rc.d/urfd/run
#!/command/with-contenv sh

# start daemon
exec /usr/local/bin/urfd ${URFD_CONFIG_DIR}/urfd.ini
